# supervisord.conf
[supervisord]
nodaemon=true
loglevel = DEBUG

[program:beanstalkd]
command=/usr/bin/beanstalkd

[program:alprd]
command=/usr/bin/alprd -f
